# sample/systemd/default.target
[Unit]
Description=default target
Wants=bash.service
After=bash.service
